# linux/systemd/rclone.timer
[Unit]
Description=Perform an rclone sync periodically

[Timer]
OnBootSec=1h
OnCalendar=02:00

[Install]
WantedBy=timers.target
